# debian/postrm
#!/bin/sh
udevadm control --reload-rules
systemctl daemon-reload || true
systemctl stop udev-mountslist.service || true
systemctl disable udev-mountslist.service || true

exit 0
